# infra/systemd/obd2-dashboard-backup-check.service
[Unit]
Description=OBD2 Dashboard Backup Integrity Check (restic)
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/obd2-dashboard-backup.env

[Service]
Type=oneshot
EnvironmentFile=/etc/obd2-dashboard-backup.env
ExecStart=/opt/obd2-dashboard/infra/scripts/check-backup.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
